# systemd/lind-ttyd.service
[Unit]
Description=Lind — ttyd terminal exposing Cline CLI
After=network.target

[Service]
Type=simple
User=exis
WorkingDirectory=/opt/lind
Environment="CLINE_API_KEY=SET_AFTER_FIRST_RUN"
Environment=PATH=/usr/bin:/usr/local/bin
ExecStart=/usr/local/bin/ttyd --port 7681 --interface 127.0.0.1 --writable /opt/lind/cline_launcher.sh
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/opt/lind
ReadOnlyPaths=/opt/lind/node_modules

[Install]
WantedBy=multi-user.target
